# load.service
[Unit]
Description=systemd service file for the load.sh script.

[Service]
ExecStart=/home/jnazirbage/load.sh

[Install]
WantedBy=multi-user.target

